# supervisord_sentry.conf
[program:sentry-web]
user=sentry
directory=/var/sentry/
command=/var/sentry/bin/sentry --config=/etc/sentry.conf.py start http
autostart=true
autorestart=true
redirect_stderr=true
