# supervisord_docker.conf
[program:sshd]
command=/usr/sbin/sshd -D

[program:mysqld]
command=service mysqld start

[program:redis]
command=service redis start
